# update-motd.d/00-figlet
#!/bin/sh

CURRENT_DIR=$(dirname "$0")
. "$CURRENT_DIR/quiet"
. "$CURRENT_DIR/colors"

printf "\n%b%s%b\n" "$LIGHT_RED" "$(figlet "  $(hostname -s)")" "$NONE"
