# update-motd.d/30-need-upgrade
#!/bin/bash

. /etc/update-motd.d/colors

n=$(printf $NEED_UPGRADE)
if [[ $n -gt 0 ]]; then
    printf $LIGHT_RED
    printf "  You have %s packages waiting for upgrades." "$n"
    printf $NONE"\n\n"
fi
n=$(printf $NOT_NEEDED)
if [[ $n -gt 0 ]]; then
    printf $YELLOW
    printf "  You have %s packages that were automatically installed and are not needed anymore." "$n"
    printf $NONE"\n\n"
fi
