# update-motd.d/40-need-reboot
#!/bin/sh

CURRENT_DIR=$(dirname "$0")
. "$CURRENT_DIR/quiet"
. "$CURRENT_DIR/colors"

if [ -e /var/run/reboot-required ]; then
    printf "%b  Pending kernel upgrade! %bYou should consider rebooting your machine.\n\n" "$LIGHT_RED" "$NONE"
fi
