# update-motd.d/colors
NONE="\e[m"
export NONE

WHITE="\e[1;37m"
export WHITE
GREEN="\e[1;32m"
export GREEN
RED="\e[0;32;31m"
export RED
YELLOW="\e[1;33m"
export YELLOW
BLUE="\e[34m"
export BLUE
CYAN="\e[36m"
export CYAN
LIGHT_GREEN="\e[1;32m"
export LIGHT_GREEN
LIGHT_RED="\e[1;31m"
export LIGHT_RED
