# update-motd.d/quiet
for i in $(cat /etc/update-motd.d/hushlogin)
do
    if [ $(cat /proc/self/loginuid) -eq $(id -u $i) ]
    then
        exit
    fi
done
